# .config/systemd/user/backup-wallpapers.service
[Unit]
Description=Backup Random Wallpaper fetched files.

[Service]
Type=oneshot
Nice=19
ExecStart=/home/phab/.local/bin/backup-wallpapers.sh
