# .local/share/systemd/user/backup-wallpapers.service
[Unit]
Description=Backup Random Wallpaper fetched files.

[Service]
Type=simple
ExecStart=/home/phab/.local/bin/backup-wallpapers.sh
