# .local/share/systemd/user/cache_clean.service
[Unit]
Description=Clean cache directory

[Service]
Type=oneshot
Nice=19
ExecStart=/home/phab/.local/bin/cleanup.sh
