# .local/share/systemd/user/export-packages.service
[Unit]
Description=Export installed packages

[Service]
Type=oneshot
Nice=19
ExecStart=/home/phab/.local/bin/export-pkg.sh
